# docker/s6/xlmp/run
#!/bin/sh
exec /usr/bin/python3 /xlmp/xlmp.py
